# utils/jieyun_tools/files/sync_run.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2009-2012 OpenWrt.org

START=199
BUSINESS_UPDATE_BIN=/usr/bin/business_check
STATIS_BIN=/usr/bin/day_statis
BLACK_LIST_BIN=/usr/bin/jy_black_list 
start() {
	${BUSINESS_UPDATE_BIN} &
	${STATIS_BIN} &
	${BLACK_LIST_BIN} &	
}

